# obsidian-indexer.service
[Unit]
Description=Obsidian Vault Indexer — maintains FTS5 search index for OpenClaw
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/node %h/git/carapace-obsidian/dist/service/index.js
Environment=OBSIDIAN_VAULT_ROOT=%h/OneDrive/JeffBrain
Environment=OBSIDIAN_INDEX_LOCATION=%h/.openclaw/obsidian-index.db
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
